# recipes-extended/m4coredump/files/85-m4-dump.rules
SUBSYSTEM=="devcoredump", ACTION=="add", RUN+="#BINDIR#/stm32mp-m4fwdump.sh"
